# mkosi.extra/etc/systemd/system/data_images.mount
[Unit]
Description=Mount the data partition
Requires=multi-user.target
After=local-fs.target
After=multi-user.target

[Mount]
What=/dev/disk/by-label/persistent_part
Where=/data_images
Type=ext4
Options=rw,noatime

[Install]
WantedBy=multi-user.target
